# btrbk-snapshot.service
[Unit]
Description=Btrbk Backup Service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/btrbk run
